# esp-idf/data/factory_state/crontab
# min hour dom mon dow flags command
# flags: - (none), A (awake only), N (network required, implies A)
*/15 * * * * N duckdns update
0 3 * * 0 N cert acme 30
