# esp-idf/data/factory_state/crontab
# crontab: periodic tasks


# flags:
#
#   -   none
#   A   When awake, not on deepsleep wakeup
#   N   When connected to on network

# min hour  dom   mon   dow flags   command
*/15  *     *     *     *     N     duckdns update
0     3     *     *     *     N     cert acme 30
