# esp-idf/data/factory_state/crontab
# crontab: periodic tasks


# flags:
#
#   -   none
#   A   When awake, not on deepsleep wakeup
#   N   When STA connected to upstream network (AP-only does not count)

# min hour  dom   mon   dow flags   command
*/15  *     *     *     *     N     duckdns update
0     3     *     *     *     N     cert acme 30
0     0     *     *     *     A     logrotate 7
